# supervisor/influxdb.conf
[program:influxdb]
directory=/data/influxdb/
command=/usr/bin/influxd -config /etc/influxdb/influxdb.conf
user=influxdb
umask=022
autostart=True
autorestart=True
redirect_stderr=True