# supervisor/influxdb.conf
[program:influxdb]
directory=/opt/influxdb/
command=/opt/influxdb/influxd -config /opt/influxdb/instance/etc/influx.conf
user=influxdb
umask=022
autostart=True
autorestart=True
redirect_stderr=True
