# supervisor/kapacitor.conf
[program:kapacitor]
directory=/data/kapacitor
command=/usr/bin/kapacitord -config /etc/kapacitor/kapacitor.conf
user=kapacitor
umask=022
autostart=True
autorestart=True
redirect_stderr=True
